# supervisord.conf.example
[program:python-slackbot]
directory=/srv/python-slackbot
command=/srv/python-slackbot/__/bin/python rtmbot.py
autostart=true
autorestart=true

